# bhoma/sysconfig/clinic_bhoma.crontab
*/4 * * * *   python /var/src/bhoma/bhoma/scripts/sshtunnel.py
*/10 * * * *  cd /var/src/bhoma/bhoma && python manage.py phonehome
0 15 * * 2    cd /var/src/bhoma/bhoma && python manage.py compact_views
