# bhoma/sysconfig/clinic_root.crontab
*/15 * * * *  /etc/init.d/modem-config start
* * * * *     cd /var/src/bhoma/bhoma && python manage.py sync_data both true
0 * * * *     /var/src/bhoma/bhoma/scripts/lsof.sh 
