# bhoma/scripts/central_server/supervisorconf.ini
[program:celery]
command=/usr/bin/python /var/src/bhoma/bhoma/manage.py celeryd --loglevel=INFO  -v 2 -s celery   
directory=/var/src/bhoma/bhoma/
user=bhoma
numprocs=1
stdout_logfile=/var/log/bhoma/celeryd.log
stderr_logfile=/var/log/bhoma/celeryd.error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
priority=998

[program:celerybeat]
command=/usr/bin/python /var/src/bhoma/bhoma/manage.py celerybeat --loglevel=INFO
directory=/var/src/bhoma/bhoma/
user=bhoma
numprocs=1
stdout_logfile=/var/log/bhoma/celerybeat.log
stderr_logfile=/var/log/bhoma/celerybeat.error.log
autostart=true
autorestart=true
startsecs=10
priority=998

